# wotstats.timer
[Unit]
Description=Update WoT stats plot daily

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
